# stacks/metapi/metapi.target
[Unit]
Description=metapi stack
Wants=metapi-pod.service tailscale-metapi.service metapi.service

[Install]
WantedBy=default.target
