# stacks/n8n/n8n.target
[Unit]
Description=n8n stack
Wants=n8n-pod.service tailscale-n8n.service n8n.service

[Install]
WantedBy=default.target
